# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:bot]
command=python -u bot.py
directory=/app
user=botuser
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/app/logs/bot.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/bot_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=
    PYTHONUNBUFFERED="1",
    TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",
    JD_EMAIL="%(ENV_JD_EMAIL)s",
    JD_PASSWORD="%(ENV_JD_PASSWORD)s",
    JD_DEVICE_NAME="%(ENV_JD_DEVICE_NAME)s",
    ALLOWED_USERS="%(ENV_ALLOWED_USERS)s",
    DOWNLOAD_PATH="%(ENV_DOWNLOAD_PATH)s",
    MAX_FILE_SIZE="%(ENV_MAX_FILE_SIZE)s",
    LOG_LEVEL="%(ENV_LOG_LEVEL)s"

[program:healthcheck]
command=python -c "from healthcheck import HealthServer; HealthServer(port=8080).start()"
directory=/app
user=botuser
autostart=true
autorestart=true
startretries=5
stdout_logfile=/app/logs/healthcheck.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
